# packages/sources/atftp-0.7.2/debian/atftpd.prerm
#!/bin/sh

set -e

# Stops daemon if it is running under our control
if [ -x /etc/init.d/atftpd ]; then 
    if [ -x /usr/sbin/invoke-rc.d ]; then
	invoke-rc.d --quiet atftpd stop
    else
	/etc/init.d/atftpd stop
    fi
fi	

#DEBHELPER#

exit 0
